# systemd/corex-nexdmarc.service
[Unit]
Description=CoreX NexDMARC
After=network.target

[Service]
Type=simple
User=nexdmarc
WorkingDirectory=/opt/corex-nexdmarc
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NEXDMARC_DATA_DIR=/var/lib/corex-nexdmarc
Environment=NEXDMARC_UPDATE_SCRIPT=/opt/corex-nexdmarc/scripts/update.sh
ExecStart=/opt/corex-nexdmarc/node_modules/.bin/tsx /opt/corex-nexdmarc/server.ts
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
